# tunasync/tunasync.ini
[program:tunasync-manager]
command=/root/tunasync/tunasync manager --config /root/tunasync/manager.conf
directory=/root/tunasync
autostart=true
autorestart=true
logfile=/var/log/supervisor/tunasync-manager.log
log_stderr=true
user=root

[program:tunasync-worker]
command=/root/tunasync/tunasync worker --config /root/tunasync/worker.conf
directory=/root/tunasync
autostart=true
autorestart=true
logfile=/var/log/supervisor/tunasync-worker.log
log_stderr=true
user=root
